# cron/crontab
*/10 * * * * python3 ~/Documents/NetworkMonitor/cron/ping_sweep.py
