# cron/crontab
*/10 * * * * python3 ~/Documents/NetworkMonitor/ping_sweep.py
