# cron/pingsweep.crontab
*/10 * * * * python3 /path/to/repo/NetworkMonitor/cron/ping_sweep.py
